# services/start_bot.service
[Unit]
Description=Start bot
[Service]
ExecStart=/usr/bin/python3 checker.py